# matrix.socket
[Socket]
ListenStream=/tmp/matrix.sk

[Install]
WantedBy=sockets.target
